# deploy/sentinel.service
[Unit]
Description=Agentic Honey Pot (Sentinel)
After=network.target

[Service]
Type=simple
User=sentinel
Group=sentinel
WorkingDirectory=/opt/agentic_honey_pot
# Set environment vars here or use an env file referenced in ExecStart
Environment=API_KEY=changeme
Environment=ADMIN_API_KEY=changeme
ExecStart=/opt/venv/bin/uvicorn backend.app.main:app --host 0.0.0.0 --port 8030
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
